# data/lidske-aktivity.timer
[Unit]
Description=Run Lidské Aktivity every 6 hours

[Timer]
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target
