# data/lidske-aktivity.service
[Unit]
Description=Lidské Aktivity

[Service]
ExecStart=/usr/bin/lidske-aktivity --scan

[Install]
WantedBy=multi-user.target
